# debian/etc/systemd/system/armcnc_python.service
[Unit]
Description=armcnc_launch
After=syslog.target network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/armcnc/python/
ExecStart=/usr/bin/python3 -B /opt/armcnc/python/launch.py
Restart=on-failure

[Install]
WantedBy=multi-user.target