# debian/DEBIAN/prerm
#!/bin/bash

# 移除之前执行业务

sudo systemctl stop armcnc.service
sudo systemctl disable armcnc.service
sudo systemctl daemon-reload